# src/mobot_pkg/services/fan.service

[Unit]
Description=PiOLED stats display service

[Service]
Type=simple
User=<user name>
ExecStart=/bin/sh -c "bash ~/fan.sh"
WorkingDirectory=/home/<user name>
Restart=always

[Install]
WantedBy=multi-user.target
